# examples/aleph/bin/solr
#!/bin/bash
#
# solr
#
# chkconfig: 345 95 25
# description: solr script
#

if [ "${SOLR_DIR}" == "" ]; then
    SOLR_DIR="./engine"
fi

if [ "${SOLR_DATA_DIR}" == "" ]; then
    SOLR_DATA_DIR="./solr/data"
fi

JAVA_OPTIONS="-Xmx1024m -Dsolr.data.dir=${SOLR_DATA_DIR} -DSTOP.PORT=8079 -DSTOP.KEY=stopkey -jar start.jar"
LOG_FILE="/tmp/solr.log"
JAVA="java"

case $1 in
    start)
        echo "Starting Solr"
        cd $SOLR_DIR
        $JAVA $JAVA_OPTIONS 2> $LOG_FILE &
        ;;
    stop)
        echo "Stopping Solr"
        cd $SOLR_DIR
        $JAVA $JAVA_OPTIONS --stop
        ;;
    status)
        status solr
        ;;
    restart)
        $0 stop
        sleep 1
        $0 start
        ;;
    *)
        echo "Usage: $0 {start|stop|restart}" >&2
        exit 1
        ;;
esac
